# gentoo/dev-lang/dm/files/dnoded.rc6
#!/sbin/runscript

[[ ${#DNODED_PORTS[@]} == 0 ]] && DNODED_PORTS=(0)
DNODED_NICE=${DNODED_NICE:-0}
DNODED_PID=${DNODED_PID:-/var/run/dnoded.pid}

depend() {
    after net xdm
}

start() {
    ebegin "Starting ${SVCNAME}"
    einfo "For ports ${DNODED_PORTS[@]}"
    start-stop-daemon -S \
	-p ${DNODED_PID} \
	-x /usr/bin/dnode-daemon \
	-- $DNODED_NICE ${DNODED_PID} ${DNODED_PORTS[@]}
    eend $?
}

stop() {
    ebegin "Stopping ${SVCNAME}"
    start-stop-daemon -K \
	-p ${DNODED_PID}
    eend $?
}
