# gentoo/dev-lang/dm/files/dnoded.rc6
#!/sbin/runscript

[[ ${#DNODED_PORTS[@]} == 0 ]] && DNODED_PORTS=(0)
DNODED_NICE=${DNODED_NICE:-0}

depend() {
    after net xdm
}

start() {
    ebegin "Starting ${SVCNAME}"
    einfo "For ports ${DNODE_PORTS[@]}"
    start-stop-daemon -S \
	-x /usr/bin/dnode-daemon \
	-- $DNODED_NICE /var/run/dnode-$i.pid ${DNODE_PORTS[@]}
    eend $?
}

stop() {
    ebegin "Stopping ${SVCNAME}"
    start-stop-daemon -K -p /var/run/dnoded-$i.pid
    eend $?
}
